# supervisord/mysql.conf
[program:mysqld]
command=/usr/sbin/mysqld --user=mysql
